# .config/systemd/user/atuin.service
[Unit]
Description=Atuin Service
After=network.target

[Service]
Type=simple
ExecStart=%h/.local/bin/atuin server start
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
