# .config/systemd/user/pueued.service
[Unit]
Description=Pueue Daemon - CLI process scheduler and manager

[Service]
Restart=no
ExecStart=%h/.local/bin/pueued -vv

[Install]
WantedBy=default.target
