# .config/systemd/user/sing-box.service
[Unit]
Description=sing-box tproxy daemon.
After=network.target nss-lookup.target network-online.target

[Service]
Type=simple
LimitNOFILE=1000000
ExecStart=%h/.local/bin/sing-box run -c %h/.config/sing-box/config.json -D %h/.config/sing-box/

[Install]
WantedBy=default.target
