# .config/systemd/user/waybar.service
[Unit]
Description=Waybar - Status bar for Wayland
Documentation=https://github.com/Alexays/Waybar
After=graphical-session.target
Wants=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStartPre=/bin/bash -c 'for i in $(seq 1 60); do busctl --user status org.freedesktop.portal.Desktop &>/dev/null && exit 0; sleep 0.5; done'
ExecStart=waybar
Restart=on-failure
RestartSec=5
Environment=PATH=%h/.local/bin:/usr/local/bin:/usr/bin
Environment=LD_LIBRARY_PATH=/usr/local/lib64:$LD_LIBRARY_PATH

[Install]
WantedBy=graphical-session.target
